# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:redis-server]
command=bach -c ". /etc/default/redis && /usr/bin/redis-server --maxclients 1000"
pidfile=/var/log/supervisor/redis-server.pid
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:apiaxle]
command=bash -c "apiaxle-proxy -f 1 -p 3000 -q "
pidfile=/var/log/supervisor/apiaxle.pid
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
